# raspi_deployer_starter/deploy/pi-backup.timer
[Unit]
Description=Daily backup

[Timer]
OnCalendar=daily
Persistent=true
Unit=pi-backup.service

[Install]
WantedBy=timers.target
